# dots/calendar/vdirsyncer.service
[Unit]
Description=Synchronize calendars and contacts
Documentation=https://vdirsyncer.readthedocs.org/

[Service]
ExecStart=/usr/local/bin/vdirsyncer sync
Type=oneshot